# falcon/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/deploy/app
command=/usr/local/bin/gunicorn -c /deploy/gunicorn_config.py api:__hug_wsgi__
autostart=true
autorestart=true
user=root
