# flask/supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
directory=/deploy/app
command=/usr/local/bin/gunicorn -c /deploy/gunicorn_config.py wsgi:app
autostart=true
autorestart=true
user=root
